# systemd/makemkv-auto-key.service
[Unit]
Description=MakeMKV Beta Key Update Service

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/makemkv-auto key update
